# services/vault/config/vault.path
[Unit]
Description=Bootstrap activator for Vault

[Path]
PathExists=/etc/consul.d/bootstrap.json

[Install]
WantedBy=default.target
