# services/consul/config/firewall.path
[Unit]
Description=firewall bootstrap

[Path]
PathExists=/etc/consul.d/bootstrap.json

[Install]
WantedBy=default.target
