# services/consul/config/firewall.path
[Unit]
Description=Bootstrap activator for the Consul firewall watcher

[Path]
PathExists=/etc/consul.d/bootstrap.json

[Install]
WantedBy=default.target
