# services/fabio/config/fabio.path
[Unit]
Description=Bootstrap activator for Fabio

[Path]
PathExists=/etc/consul.d/bootstrap.json

[Install]
WantedBy=default.target
