# services/consul/config/telegraf/telegraf.path
[Unit]
Description=Bootstrap activator for Telegraf

[Path]
PathExists=/etc/consul.d/bootstrap.json

[Install]
WantedBy=default.target
