# services/influx/config/chronograf/chronograf.path
[Unit]
Description=Bootstrap activator for Chronograf

[Path]
PathExists=/mnt/persist/chronograf/

[Install]
WantedBy=default.target
