# services/influx/config/influx.path
[Unit]
Description=Bootstrap activator for Influx

[Path]
PathExists=/mnt/persist/influx/

[Install]
WantedBy=default.target
