# images/consul/config/consul.path
[Unit]
Description=Bootstrap activator for Consul

[Path]
PathExists=/etc/consul.d/bootstrap.json

[Install]
WantedBy=default.target
