# services/nomad/config/nomad.path
[Unit]
Description=nomad bootstrap

[Path]
PathExists=/etc/nomad.d/bootstrap.hcl

[Install]
WantedBy=default.target
